# deploy/invoice-collector.service
[Unit]
Description=Invoice Collector Service
After=network.target postgresql.service

[Service]
Type=simple
User=tran-ninh
Group=tran-ninh
WorkingDirectory=/home/tran-ninh/OtherProjects/invoice_server
Environment=PYTHONPATH=/home/tran-ninh/OtherProjects/invoice_server
ExecStart=/usr/bin/python3 -m backend.collector.main
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=invoice-collector

[Install]
WantedBy=multi-user.target
